# package/src/DEBIAN/postinst
#!/bin/bash
mv /etc/poseidonos.service /lib/systemd/system/poseidonos.service
systemctl daemon-reload
systemctl enable poseidonos.service
systemctl start poseidonos.service
