# package/src/DEBIAN/postinst
#!/bin/bash

# libtcmalloc linking
ln -sf /usr/local/lib/libtcmalloc.so.4.5.3 /usr/local/lib/libtcmalloc.so.4
ln -sf /usr/local/lib/libtcmalloc.so.4 /usr/local/lib/libtcmalloc.so

# Register CLI man page
mv /etc/poseidonos-cli* /usr/share/man/man3/
makewhatis

# register poseidonos.service
mv /etc/poseidonos.service /lib/systemd/system/poseidonos.service
systemctl daemon-reload
systemctl enable poseidonos.service
