# package/etc/poseidonos.service
[Unit]
Description=poseidonos

[Service]
ExecStart=/usr/local/bin/poseidonos

[Install]
WantedBy=multi-user.target
Alias=poseidonos
